# cloud_sync/crontab
*/5 * * * * /usr/bin/env python /app/cloud_sync/cloud_sync.py
